# overwatcher-dashboard.service
[Unit]
Description=OverwatcherPI Streamlit Dashboard
After=network.target

[Service]
Type=simple
User=root
# Adjust User and WorkingDirectory to match your deployment path (usually /opt/OverwatcherPI for prod)
WorkingDirectory=/opt/OverwatcherPI
Environment="PATH=/opt/OverwatcherPI/dashboard/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/OverwatcherPI/dashboard/venv/bin/streamlit run /opt/OverwatcherPI/dashboard/app.py --server.address 0.0.0.0 --server.port 8109 --server.headless true
Restart=always
RestartSec=5
StandardOutput=append:/opt/OverwatcherPI/dashboard/dashboard.log
StandardError=append:/opt/OverwatcherPI/dashboard/dashboard.log

[Install]
WantedBy=multi-user.target
