# overwatcher-dashboard.service
[Unit]
Description=OverwatcherPI Streamlit Dashboard
After=network.target

[Service]
Type=simple
User=abhipra
WorkingDirectory=/mnt/usb-storage/Work/OverwatcherPI
Environment="PATH=/mnt/usb-storage/Work/OverwatcherPI/dashboard/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/mnt/usb-storage/Work/OverwatcherPI/dashboard/venv/bin/streamlit run /mnt/usb-storage/Work/OverwatcherPI/dashboard/app.py --server.address 127.0.0.1 --server.port 8501 --server.headless true
Restart=always
RestartSec=5
StandardOutput=append:/mnt/usb-storage/Work/OverwatcherPI/dashboard/dashboard.log
StandardError=append:/mnt/usb-storage/Work/OverwatcherPI/dashboard/dashboard.log

[Install]
WantedBy=multi-user.target
